# addon/qmail-init-1.03/etc/profile.d/qmail.sh
MAIL=${HOME}/Maildir/
MAILDROP="$MAIL"
export MAIL MAILDROP

test -d Maildir || maildirmake Maildir
